# misc/98-usb-gsm.rules
ATTRS{idVendor}=="12d1", ATTRS{idProduct}=="1446", RUN+="/root/modeswitch_e173.sh"
